# supervisord.conf
# =============================================================================
#  OpenMap — supervisord configuration
#  NOTE: This file is regenerated at container start by init.sh.
#        Optional services (photon, osrm) are only added when their data
#        is available. This file serves as the template/fallback.
# =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/openmap/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/openmap/supervisord.pid
user=root

[unix_http_server]
file=/var/run/openmap/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/openmap/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ── Nginx (always on) ────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

# ── Photon (optional — enabled by init.sh if data exists) ───────────────────
[program:photon]
command=java -Xmx4g -jar /opt/photon/photon.jar -data-path /data/photon -listen-port 2322
autostart=false
autorestart=true
startretries=3
startsecs=15
priority=20
stdout_logfile=/var/log/openmap/photon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/openmap/photon-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2

# ── OSRM (optional — enabled by init.sh if data exists) ─────────────────────
[program:osrm]
command=osrm-routed --algorithm=MLD /data/osrm/germany.osrm --port 5000 --max-table-size 10000
autostart=false
autorestart=true
startretries=3
startsecs=5
priority=20
stdout_logfile=/var/log/openmap/osrm.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/var/log/openmap/osrm-err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
